# goldenstatue/data/supervisor/goldenstatue_static.conf
[program:goldenstatue_static]
directory = /var/www/BakeWebOut/goldenstatue/static
environment = NODE_ENV="production"
command = node src2/server.js
user = www-data
autostart = true
autorestart = true
startsecs = 5
redirect_stderr = true
stdout_logfile = /var/www/BakeWebOut/goldenstatue/data/log/static-std.log