# goldenstatue/data/supervisor/goldenstatue_server.conf
[program:goldenstatue_server]
directory = /var/www/BakeWebOut/goldenstatue/server
command = /var/www/BakeWebOut/goldenstatue/server/goldenstatue
user = www-data
autostart = true
autorestart = true
startsecs = 5
redirect_stderr = true
stdout_logfile = /var/www/BakeWebOut/goldenstatue/data/log/server-std.log

